# scripts/local.d/cpufreq.start
#!/bin/sh

cpupower frequency-set -g ondemand
