# scripts/local.d/kvmsupp.start
#! /bin/sh

grep "^flags[[:space:]]\+: " /proc/cpuinfo | grep -qw svm && modprobe -q kvm-amd
grep "^flags[[:space:]]\+: " /proc/cpuinfo | grep -qw vmx && modprobe -q kvm-intel
