# scripts/local.d/media.start
#! /bin/sh

findmnt /media > /dev/null || mount -t tmpfs tmpfs /media
