# scripts/local.d/net.stop
#! /bin/sh

nic=wlan0
[ ! -e /sys/class/net/$nic ] && exit 0
#dhcpcd -x [iface]
start-stop-daemon -K -q -p /var/run/dhcpcd-$nic.pid > /dev/null 2>&1
start-stop-daemon -K -q -p /var/run/wpa_supplicant.pid > /dev/null 2>&1
rm -f /var/run/wpa_supplicant.pid
ip link set dev $nic down
