# scripts/local.d/sound.stop
#! /bin/sh

#start-stop-daemon -K -q -p /var/run/start_sound.pid
start-stop-daemon -K -q -p /var/run/start_sound.pid
